# telescope.service
[Unit]
Description=Telescope HUD (Cage)
After=multi-user.target

[Service]
# Run as root to ensure full access to DRM/KMS and TTY
User=root
Group=root

# Essential Wayland environment variables
RuntimeDirectory=telescope
Environment=XDG_RUNTIME_DIR=/run/telescope
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=WLR_BACKENDS=drm

# Launch Cage in "scene" mode (-s) using the setup script
ExecStart=/usr/bin/cage -s -- /home/raspi/telescope-raspi/run_cage.sh
WorkingDirectory=/home/raspi/telescope-raspi

# Log output to a file for easier debugging
StandardOutput=append:/var/log/telescope.log
StandardError=inherit

Restart=on-failure
RestartSec=5

# Required for hardware/GPU access from a service
PAMName=login
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
